# build/assets/firstboot/run
#!/sbin/openrc-run

command="/usr/bin/first-boot"
command_background=false
depend() {
  after modules
  need localmount
}
