# build/assets/pi-init/run
#!/sbin/openrc-run

command="/usr/bin/pi-init"
command_background=false
depend() {
  after modules
  need localmount
}
